# systemd/himawari.service
# $HOME/.config/user/systemd/himawari.service

[Unit]
Description=Update desktop background with satellite imagery

[Service]
Type=oneshot
ExecStart=<INSTALLATION_PATH>
